# app/crontab
1 * 4,24 * * cd /usr/src/app && python app.py 
